# Final/Product/website.service
[Unit]
Description=webserver

[Service]
ExecStart=/usr/bin/java -jar /www/servers-shop-0.0.1-SNAPSHOT.jar

[Install]
WantedBy=multi-user.target